# deploy/haruhifanclub-test.service
# 临时测试后端 systemd 单元（test.haruyuki.cn 专用，独立于生产）。
# 安装：cp 到 /etc/systemd/system/ ; systemctl daemon-reload ; systemctl enable --now haruhifanclub-test
[Unit]
Description=haruhifanclub TEST backend (axum) for test.haruyuki.cn
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/haruhifanclub-test
EnvironmentFile=/var/www/haruhifanclub-test/.env
ExecStart=/var/www/haruhifanclub-test/bin/haruhi-server
Restart=on-failure
RestartSec=3
# 给优雅关闭留时间（刷 WAL + 处理在途请求），超时才 SIGKILL
KillSignal=SIGTERM
TimeoutStopSec=30
# 数据与上传目录需可写
ReadWritePaths=/var/www/haruhifanclub-test/data /var/www/haruhifanclub-test/uploads

# ---- 内存约束（与生产对等，详见 haruhifanclub.service）----
MemoryHigh=256M
MemoryMax=512M
Environment=MALLOC_ARENA_MAX=2

# ---- 沙箱加固（与生产对等）----
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
UMask=0077

[Install]
WantedBy=multi-user.target
